# examples/zookeeper/pod.conf
[program:pod]
command=python /opt/pod/pod.py
